# src/ssh/supervisor.ssh.conf
[program:sshd]
process_name = sshd
command = /usr/sbin/sshd -D
autostart = true
autorestart = true
